# templates/default/upstart/eventstore.conf.erb
pre-start script

bash << "EOF"
  mkdir -p <%= @logs_dir %>
  chown -R <%= @user %> <%= @logs_dir %>
  mkdir -p <%= @data_dir %>
  chown -R <%= @user %> <%= @data_dir %>
EOF

end script

start on (started network-interface
          or started network-manager
          or started networking
          or runlevel [2345])

stop on (stopping network-interface
         or stopping network-manager
         or stopping networking
         or runlevel [!2345])

respawn

env MONO_GC_DEBUG=clear-at-gc
env LD_LIBRARY_PATH=<%= @executable_dir %>
env MONO_THREADS_PER_CPU=100

exec start-stop-daemon --start --chuid <%= @user %> --chdir <%= @executable_dir %> --exec ./run-node.sh -- --config=<%= @config_file %>
